# FrontEnd/service/bin/run
#!/usr/bin/env node
// Load package.json
const pjs = require('../package.json');
// Get some meta info from the package.json
const { name, version } = pjs;

const ip = require('ip');
const SERVICE_NAME = name;
const SCHEME = 'http';
const HOST = ip.address();
const PORT = 5050;

const consul = require('consul')({host: '172.16.238.2', port:8500});
const http = require('http');

const config = require('../config')[process.env.NODE_ENV || 'development'];
const log = config.log();
const SERVICE_ID = config.serviceID;
const service = require('../server/service')(config);
const server = http.createServer(service);

// Important - a service should not have a fixed port but should randomly choose one
server.listen(PORT);

server.on('listening', () => {
  var check = {
	  id: SERVICE_ID,
	  name: SERVICE_NAME,
	  address: HOST,
	  port: PORT,
	  check: {
		http: SCHEME+'://'+HOST+':'+PORT+'/frontEnd/prove',
		ttl: '5s',
		interval: '3s',
		timeout: '3s',
		deregistercriticalserviceafter: '3s'
	  }
  };
  consul.agent.service.register(check, function(err) {
		if (err) throw err;
  });
  
  log.info(
    `Hi there! I'm listening on IP ${ip.address()} on port ${server.address().port} in ${service.get('env')} mode.`
  );
  
});
